# setups/waybar/systemd/waybar.service
[Unit]
Description=Waybar status bar

[Service]
Type=simple
ExecStart=/usr/bin/waybar
Restart=on-failure

[Install]
WantedBy=default.target


