# setups/mako/systemd/mako.service
[Unit]
Description=Mako notification daemon

[Service]
Type=simple
ExecStart=/usr/bin/mako
Restart=on-failure

[Install]
WantedBy=default.target


